# automation/com.dfiore.ai-hedge-fund.morning-brief.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.dfiore.ai-hedge-fund.morning-brief</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd "/Users/davidfiore/Documents/Hedge Fund/current-ai-hedge-fund" &amp;&amp; mkdir -p reports/morning_brief &amp;&amp; .venv/bin/python main.py morning-email today &gt;&gt; reports/morning_brief/automation.log 2&gt;&amp;1</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/davidfiore/Documents/Hedge Fund/current-ai-hedge-fund</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>4</integer>
    <key>Minute</key>
    <integer>45</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/davidfiore/Documents/Hedge Fund/current-ai-hedge-fund/reports/morning_brief/launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/davidfiore/Documents/Hedge Fund/current-ai-hedge-fund/reports/morning_brief/launchd.err.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
