# automation/com.dfiore.ai-hedge-fund.daily-setup-review.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.dfiore.ai-hedge-fund.daily-setup-review</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/davidfiore/Documents/Hedge Fund/current-ai-hedge-fund/scripts/run_daily_setup_review.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/davidfiore/Documents/Hedge Fund/current-ai-hedge-fund</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>16</integer>
    <key>Minute</key>
    <integer>20</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/davidfiore/Documents/Hedge Fund/current-ai-hedge-fund/reports/setup_review/launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/davidfiore/Documents/Hedge Fund/current-ai-hedge-fund/reports/setup_review/launchd.err.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
